# opencode/cron/crontab
# OpenCode Scheduled Tasks (opencode user)
# Edit this file, then: docker compose restart opencode-daemon

SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
HOME=/home/opencode
MAILTO=""
WORKSPACE=/workspace
cd /workspace

# Format: minute hour day month weekday command
# Examples: 0 23 * * * (11 PM daily), 0 8 * * * (8 AM), */30 * * * * (every 30 min), 0 9 * * 1 (Mon 9 AM)

# Nightly reflection - proposals to SOUL.md.proposals (11 PM)
0 23 * * * /workspace/scripts/run-task.sh "Perform nightly reflection: Review today's activities from logs. Reflect on what went well, what could improve, and any patterns noticed. Append your proposals to /workspace/SOUL.md.proposals using the format documented in AGENTS.md. Do NOT modify SOUL.md directly." reflection

# Add your jobs below (or uncomment examples):
# 0 8 * * * /workspace/scripts/run-task.sh "Morning check: Verify TrueNAS is healthy, check arr stack imports." morning-check
# 0 */6 * * * /workspace/scripts/run-task.sh "Check Discord for mentions and respond if needed." discord-check
